# tunnelManager.conf
description "tunnelManager"

start on runlevel [2345]

stop on runlevel [016]

respawn
respawn limit 10 5



chdir /data/platform_scripts/tunnelManager/tunnelManager

exec /usr/bin/python tunnelManager.py 2>&1 | tee -a /var/log/capillary/tunnelManager/tunnelmanager.out

